# deploy/devpane.service
[Unit]
Description=DevPane autonomous development daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/devpane
ExecStartPre=/usr/bin/env pnpm install --frozen-lockfile
ExecStartPre=/usr/bin/env pnpm build
ExecStart=/usr/bin/env node packages/daemon/dist/index.js
Environment=NODE_ENV=production
EnvironmentFile=-/opt/devpane/.env
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
LimitNOFILE=65536

[Install]
WantedBy=default.target
